# rules.d/1.3_files_dirs.rules
## Privilege Abuse
### The purpose of this rule is to detect when an admin may be abusing power by looking in user's home dir.
-a always,exit -F arch=b64 -F dir=/home -F uid=0 -F auid>=1000 -F auid!=-1 -C auid!=obj_uid -k T1078.003.Whome_dirs

# Audit Configuration Monitoring
-a always,exit -F arch=b64 -F path=/etc/libaudit.conf -F perm=wa -k T1562.001.Wlibaudit_conf
-a always,exit -F arch=b64 -F dir=/etc/audit/ -F perm=wa -k T1562.001.Waudit_dir
-a always,exit -F arch=b64 -F dir=/var/log/audit/ -F perm=rwa -k T1070.Waudit_logs

# Kernel Configuration Monitoring
-a always,exit -F arch=b64 -F dir=/etc/modprobe.d/ -F perm=wa -k T1542.001.Wmodprobe_d
-a always,exit -F arch=b64 -F dir=/etc/sysctl.d/ -F perm=wa -k T1542.Wsysctl_d

# Time Configuration Monitoring
-a always,exit -F arch=b64 -F path=/etc/localtime -F perm=wa -k T1070.004.Wlocaltime

# Cron and Scheduled Jobs Monitoring
-a always,exit -F arch=b64 -F path=/etc/cron.allow -F perm=wa -k T1053.003.Wcron_allow
-a always,exit -F arch=b64 -F path=/etc/cron.deny -F perm=wa -k T1053.003.Wcron_deny
-a always,exit -F arch=b64 -F path=/etc/crontab -F perm=wa -k T1053.003.Wcrontab
-a always,exit -F arch=b64 -F path=/etc/inittab -F perm=wa -k T1543.001.Winittab
-a always,exit -F arch=b64 -F path=/etc/at.allow -F perm=wa -k T1053.001.Wat_allow
-a always,exit -F arch=b64 -F path=/etc/at.deny -F perm=wa -k T1053.001.Wat_deny
-a always,exit -F arch=b64 -F path=/etc/anacrontab -F perm=wa -k T1053.003.Wanacrontab
-a always,exit -F arch=b64 -F dir=/etc/cron.d/ -F perm=wa -k T1053.003.Wcron_d
-a always,exit -F arch=b64 -F dir=/etc/cron.daily/ -F perm=wa -k T1053.003.Wcron_daily
-a always,exit -F arch=b64 -F dir=/etc/cron.hourly/ -F perm=wa -k T1053.003.Wcron_hourly
-a always,exit -F arch=b64 -F dir=/etc/cron.monthly/ -F perm=wa -k T1053.003.Wcron_monthly
-a always,exit -F arch=b64 -F dir=/etc/cron.weekly/ -F perm=wa -k T1053.003.Wcron_weekly
-a always,exit -F arch=b64 -F dir=/etc/init.d/ -F perm=wa -k T1543.001.Winit_d
-a always,exit -F arch=b64 -F dir=/etc/init/ -F perm=wa -k T1543.001.Winit_dir
-a always,exit -F arch=b64 -F dir=/var/spool/at/ -F perm=wa -k T1053.001.Wspool_at
-a always,exit -F arch=b64 -F dir=/var/spool/cron/ -F perm=wa -k T1053.003.Wspool_cron

# Account Management Monitoring
-a always,exit -F arch=b64 -F path=/etc/sudoers -F perm=rwa -k T1548.003.Wsudoers

# Session Monitoring
-a always,exit -F arch=b64 -F path=/var/run/utmp -F perm=wa -k T1070.001.Wutmp
-a always,exit -F arch=b64 -F path=/var/log/wtmp -F perm=wa -k T1070.001.Wwtmp
-a always,exit -F arch=b64 -F path=/var/log/btmp -F perm=wa -k T1070.001.Wbtmp

# Login Monitoring
-a always,exit -F arch=b64 -F path=/var/log/faillog -F perm=wa -k T1110.Wfaillog
-a always,exit -F arch=b64 -F path=/var/log/lastlog -F perm=wa -k T1070.001.Wlastlog
-a always,exit -F arch=b64 -F path=/var/log/tallylog -F perm=wa -k T1110.Wtallylog

# PAM Configuration Monitoring
-a always,exit -F arch=b64 -F path=/etc/security/pam_env.conf -F perm=wa -k T1556.Wpam_env
-a always,exit -F arch=b64 -F path=/etc/security/namespace.conf -F perm=wa -k T1556.Wnamespace_conf
-a always,exit -F arch=b64 -F path=/etc/security/namespace.init -F perm=wa -k T1556.Wnamespace_init
-a always,exit -F arch=b64 -F path=/etc/security/limits.conf -F perm=wa -k T1556.Wlimits_conf
-a always,exit -F arch=b64 -F path=/etc/pam.d/common-password -F perm=wa -k T1556.Wcommon_password
-a always,exit -F arch=b64 -F dir=/etc/pam.d/ -F perm=wa -k T1556.Wpam_d
-a always,exit -F arch=b64 -F dir=/etc/security/limits.d -F perm=wa -k T1556.Wlimits_d
-a always,exit -F arch=b64 -F dir=/etc/security/namespace.d -F perm=wa -k T1556.Wnamespace_d

# SSH Configuration Monitoring
-a always,exit -F arch=b64 -F path=/etc/ssh/sshd_config -F perm=wa -k T1021.004.Wsshd_config
-a always,exit -F arch=b64 -F dir=/etc/ssh/sshd_config.d/ -F perm=wa -k T1021.004.Wsshd_config_d

# Privilege Escalation Monitoring
-a always,exit -F arch=b64 -F dir=/usr/bin/ -F perm=wa -k T1543.Wusr_bin

## Recon Related Events
### Read
-a always,exit -F arch=b64 -F path=/etc/locale.conf -F perm=r -k T1083.Rlocale_conf
-a always,exit -F arch=b64 -F path=/etc/hostname -F perm=r -k T1083.Rhostname
-a always,exit -F arch=b64 -F path=/etc/os-release -F perm=r -k T1083.Ros_release
-a always,exit -F arch=b64 -F path=/etc/issue -F perm=r -k T1083.Rissue
-a always,exit -F arch=b64 -F path=/etc/group -F perm=r -k T1083.Rgroup
-a always,exit -F arch=b64 -F path=/etc/gshadow -F perm=r -k T1083.Rgshadow
-a always,exit -F arch=b64 -F path=/etc/shadow -F perm=r -k T1003.Rshadow
-a always,exit -F arch=b64 -F path=/etc/security/opasswd -F perm=r -k T1003.Ropasswd
-a always,exit -F arch=b64 -F path=/etc/passwd -F perm=r -k T1083.Rpasswd
-a always,exit -F arch=b64 -F path=/etc/login.defs -F perm=r -k T1083.Rlogin_defs
-a always,exit -F arch=b64 -F path=/etc/securetty -F perm=r -k T1083.Rsecuretty
-a always,exit -F arch=b64 -F path=/etc/hosts -F perm=r -k T1016.Rhosts
-a always,exit -F arch=b64 -F path=/etc/networks -F perm=r -k T1016.Rnetworks

# Credential File Modifications (Write Operations)
-a always,exit -F arch=b64 -F path=/etc/group -F perm=wa -k T1098.Wgroup
-a always,exit -F arch=b64 -F path=/etc/gshadow -F perm=wa -k T1098.Wgshadow
-a always,exit -F arch=b64 -F path=/etc/shadow -F perm=wa -k T1098.Wshadow
-a always,exit -F arch=b64 -F path=/etc/security/opasswd -F perm=wa -k T1098.Wopasswd
-a always,exit -F arch=b64 -F path=/etc/passwd -F perm=wa -k T1098.Wpasswd
-a always,exit -F arch=b64 -F path=/etc/login.defs -F perm=wa -k T1098.Wlogin_defs
-a always,exit -F arch=b64 -F path=/etc/securetty -F perm=wa -k T1098.Wsecuretty
-a always,exit -F arch=b64 -F path=/etc/hosts -F perm=wa -k T1562.001.Whosts
-a always,exit -F arch=b64 -F path=/etc/networks -F perm=wa -k T1562.001.Wnetworks

# Package Manager Monitoring
-a always,exit -F arch=b64 -F path=/etc/dnf/dnf.conf -F perm=wa -k T1543.001.Wdnf_conf
-a always,exit -F arch=b64 -F dir=/etc/yum.repos.d/ -F perm=wa -k T1543.001.Wyum_repos
-a always,exit -F arch=b64 -F path=/etc/apt/sources.list -F perm=wa -k T1543.001.Wapt_sources
-a always,exit -F arch=b64 -F dir=/etc/apt/sources.list.d/ -F perm=wa -k T1543.001.Wapt_sources_d

# Security Module Monitoring
-a always,exit -F arch=b64 -F path=/etc/selinux/config -F perm=wa -k T1562.001.Wselinux_config
-a always,exit -F arch=b64 -F dir=/etc/selinux/ -F perm=wa -k T1562.001.Wselinux_dir
-a always,exit -F arch=b64 -F dir=/etc/apparmor/ -F perm=wa -k T1562.001.Wapparmor_dir

# Library Path Monitoring
-a always,exit -F arch=b64 -F path=/etc/ld.so.conf -F perm=a -k T1574.001.Wld_so_conf
-a always,exit -F arch=b64 -F dir=/etc/ld.so.conf.d/ -F perm=wa -k T1574.001.Wld_so_conf_d

# SSH Key Monitoring
-a always,exit -F arch=b64 -F dir=/root/.ssh/ -F perm=wa -k T1098.004.Wroot_ssh

# Systemd Monitoring
-a always,exit -F arch=b64 -F dir=/usr/lib/systemd/system-generators/ -F perm=wa -k T1543.002.Wsystemd_generators
-a always,exit -F arch=b64 -F dir=/lib/systemd/system-generators/ -F perm=wa -k T1543.002.Wsystemd_generators
-a always,exit -F arch=b64 -F dir=/etc/systemd/ -F perm=wa -k T1543.002.Wsystemd_conf
-a always,exit -F arch=b64 -F dir=/usr/lib/systemd -F perm=wa -k T1543.002.Wsystemd_lib

# Shell Configuration Monitoring
-a always,exit -F arch=b64 -F path=/etc/profile -F perm=wa -k T1546.004.Wprofile
-a always,exit -F arch=b64 -F path=/etc/shells -F perm=wa -k T1546.004.Wshells
-a always,exit -F arch=b64 -F path=/etc/bashrc -F perm=wa -k T1546.004.Wbashrc
-a always,exit -F arch=b64 -F path=/etc/csh.cshrc -F perm=wa -k T1546.004.Wcshrc
-a always,exit -F arch=b64 -F path=/etc/csh.login -F perm=wa -k T1546.004.Wcsh_login
-a always,exit -F arch=b64 -F path=/etc/fish -F perm=wa -k T1546.004.Wfish_config
-a always,exit -F arch=b64 -F path=/etc/zsh -F perm=wa -k T1546.004.Wzsh_config
-a always,exit -F arch=b64 -F dir=/etc/profile.d/ -F perm=wa -k T1546.004.Wprofile_d

# Docker Monitoring
-a always,exit -F arch=b64 -F path=/usr/lib/systemd/system/docker.service -F perm=wa -k T1613.Wdocker_service
-a always,exit -F arch=b64 -F path=/usr/lib/systemd/system/docker.socket -F perm=wa -k T1613.Wdocker_socket
-a always,exit -F arch=b64 -F dir=/var/lib/docker -F perm=wa -k T1613.Wdocker_lib
-a always,exit -F arch=b64 -F dir=/etc/docker -F perm=wa -k T1613.Wdocker_conf

# Critical Directory Monitoring
-a always,exit -F arch=b64 -F dir=/tmp -F perm=wa -k T1070.004.Wtmp_dir
-a always,exit -F arch=b64 -F dir=/var/tmp -F perm=wa -k T1070.004.Wvar_tmp
-a always,exit -F arch=b64 -F dir=/usr -F perm=wa -k T1070.004.Wusr_dir
-a always,exit -F arch=b64 -F dir=/root -F perm=wa -k T1070.004.Wroot_dir
-a always,exit -F arch=b64 -F dir=/home -F perm=wa -k T1070.004.Whome_dir
-a always,exit -F arch=b64 -F dir=/opt -F perm=wa -k T1070.004.Wopt_dir

